# etc/systemd/system/steve.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/steve --jobs=16 --load-average=17
