# etc/systemd/system/nginx.service.d/override.conf
[Unit]
After=network-online.target
Requires=network-online.target
